# builder/data/32bit/etc/systemd/system/bthelper@.service.d/override.conf
[Unit]
After=hciuart.service bluetooth.service
Before=

[Service]
ExecStartPre=/bin/sleep 5